# scripts/rpi-eink/trainboard-display.service
[Unit]
Description=Trainboard E-Ink Display Update Service
After=network.target

[Service]
Type=oneshot
User=pi
Group=pi
Environment="HOME=/home/pi"
Environment="APP_DIR=/home/pi/trainboard"
Environment="DISPLAY_DIR=/home/pi/e-Paper/RaspberryPi_JetsonNano/python"
ExecStart=/home/pi/trainboard/scripts/rpi-eink/update-display.sh
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=true

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
